# private_dot_config/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Exec="flameshot &"
Type=Application
